# etc/supervisor/proc/apps.conf
[program:web]
command=node ../../web/app.js
startretries=10000
autostart=true
autorestart=true
serverurl=AUTO
stderr_logfile=../../log/web.err.log
stdout_logfile=../../log/web.out.log
