# etc/supervisor/proc/apps.conf
[program:api]
command=node ../../api/app.js
startretries=10000
autostart=true
autorestart=true
serverurl=AUTO
stderr_logfile=../../log/api.err.log
stdout_logfile=../../log/api.out.log
